# devOps/webhook/webhook.conf
# /etc/init/webhook.conf

description "webhook"
start on stopped rc RUNLEVEL=[2345]
respawn
exec python3 /home/ec2-user/app/main.py
